# udev/rules.d/11-sda-apm-fix.rules
ACTION=="add", SUBSYSTEM=="block", KERNEL=="sda", RUN+="/usr/bin/hdparm -B254  /dev/sda"
