# profile.d/70-systemd-shell-extra.sh
../../usr/lib/systemd/profile.d/70-systemd-shell-extra.sh