# profile.d/80-systemd-osc-context.sh
../../usr/lib/systemd/profile.d/80-systemd-osc-context.sh