# profile.d/cxoffice.csh
[ -d /opt/cxoffice/bin ] && setenv PATH ${PATH}:/opt/cxoffice/bin
