# profile.d/cxoffice.sh
[ -d /opt/cxoffice/bin ] && export PATH="${PATH}:/opt/cxoffice/bin"
