# systemd/system/myscripts.service
[Unit]
Description=My script

[Service]
ExecStart=/home/mojtaba/Scripts/my-script

[Install]
WantedBy=multi-user.target
 
