# systemd/system/polipo.service
.include /usr/lib/systemd/system/polipo.service
[Service]
User=polipo

