# systemd/system/timer-daily.target.wants/modprobed_db.service
[Unit]
Description=Run modprobed_db

[Service]
User=mojtaba
ExecStart=/usr/bin/modprobed_db store
