# xdg/autostart/gnome-shell-overrides-migration.desktop
[Desktop Entry]
Type=Application
Name=GNOME settings overrides migration
NoDisplay=true
Exec=/usr/lib/gnome-shell-overrides-migration.sh
